# supervisor/config/pac-worker.ini
[program:pac-admin-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pac/pac_admin/artisan queue:work --queue=high,default,low --once --sleep=3 --tries=3 --delay=30
numprocs=8
priority=999
autostart=true
autorestart=true
startsecs=1
startretries=3
stdout_logfile=/var/www/pac/pac_admin/storage/logs/worker.log

[program:pac-user-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pac/pac_user_api/artisan queue:work --queue=high,default,low --once --sleep=3 --tries=3 --delay=30
numprocs=8
priority=999
autostart=true
autorestart=true
startsecs=1
startretries=3
stdout_logfile=/var/www/pac/pac_user_api/storage/logs/worker.log